# etc/supervisord.d/jenkins.ini
[program:jenkins-master]
command=/opt/qnib/jenkins/master/bin/start.sh
startsecs=5
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
